# script/vmevents.service
[Unit]
Description=Executes getevents.sh

[Service]
Type=oneshot
ExecStart=/opt/vmevents/getevents.sh