# script/vmevents.service
[Unit]
Description=Executes getevents.sh

[Service]
Type=simple
ExecStart=/opt/vmevents/getevents.sh