# script/vmevents.timer
[Unit]
Description=Run vmevents.service every 5 minutes

[Timer]
OnCalendar=*:0/05